# hadoop-cdh5-mapreduce-historyserver/files/supervisord.d/hdfs-namenode.conf
[supervisord]
nodaemon=true
[program:hadoop-hdfs-namenode]
autostart = True
autorestart = True
command=/etc/init.d/hadoop-hdfs-namenode start
stdout_logfile = /var/log/hadoop-hdfs/supervisord.stdout.log
stderr_logfile = /var/log/hadoop-hdfs/supervisord.stderr.log
stdout_logfile_maxbytes = 10MB
stderr_logfile_maxbytes = 10MB
